# udev/80-keychron.rules
SUBSYSTEMS=="usb", ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="0250", ACTION=="add", RUN+="/bin/bash /home/olback/.dotfiles/bin/fix-keychron.sh"

